# scanner.crontab
* * * * * curl http://localhost:8000/video/scanner
